# deploy/playbooks/roles/common/templates/systemd/chacra.service.j2
# {{ ansible_managed }}
[Unit]
Description=chacra gunicorn service
After=network.target chacra-celery.service
Requires=chacra-celery.service

[Service]
Type=simple
ExecStart={{ app_home }}/bin/gunicorn_pecan -w 10 -t 300 {{ app_home }}/src/{{ app_name }}/prod.py
EnvironmentFile=-/etc/conf.d/{{ app_name }}-sysconfig
User={{ ansible_ssh_user }}
WorkingDirectory={{ app_home }}/src/{{ app_name }}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
